# pkg/freebsd/netmagis-topo/pkg-plist
etc/rc.d/topographd.sh
etc/rc.d/toposendd.sh
lib/netmagis/topo/anaconf
lib/netmagis/topo/analyze
lib/netmagis/topo/buildgraph
lib/netmagis/topo/cache
lib/netmagis/topo/checkdns
lib/netmagis/topo/decodename
lib/netmagis/topo/doteq
lib/netmagis/topo/dumpgraph
lib/netmagis/topo/extractcoll
lib/netmagis/topo/extracteq
lib/netmagis/topo/extractl1
lib/netmagis/topo/extractl2
lib/netmagis/topo/extractl3
lib/netmagis/topo/getnetif
lib/netmagis/topo/list-rancid
lib/netmagis/topo/list-vlans
lib/netmagis/topo/parse-cisco.tcl
lib/netmagis/topo/parse-hp.tcl
lib/netmagis/topo/parse-juniper.tcl
lib/netmagis/topo/parse-server.tcl
lib/netmagis/topo/sendmail
lib/netmagis/topo/ssidsensors
lib/netmagis/topo/start-rancid
lib/netmagis/topo/topographd
lib/netmagis/topo/toposendd
lib/netmagis/topo/updatecache
@dirrmtry var/netmagis/eqvirt
@dirrmtry var/netmagis/eqgen
@dirrmtry var/netmagis/cache
@dirrmtry var/netmagis
@dirrmtry var
@dirrmtry lib/netmagis/topo
@dirrmtry lib/netmagis
@exec mkdir -p %D/var/netmagis/eqvirt
@exec mkdir -p %D/var/netmagis/eqgen
@exec mkdir -p %D/var/netmagis/cache
